# package/DEBIAN/postinst
rm /var/www/html/index.html
chown -R www-data:www-data /var/www/html
chmod -R 0755 /var/www/html
install -o root -g root -m 0755 /tmp/raprec.service /lib/systemd/system
[ -f /etc/default/raprec ] || install -o root -g root -m 0644 -T /tmp/raprec.defaults /etc/default/raprec
systemctl daemon-reload
systemctl enable raprec
cp /tmp/hostapd-rap.conf /etc/hostapd
cp /tmp/dnsmasq.conf /etc
mkdir /var/www/html/vids
chown -R www-data:www-data /var/www/html

#Enable raspi camera module
raspi-config nonint do_camera 1

# Enable SSL and HTTP => HTTPS redirect.
install -o root -g root -m 0644 /tmp/rewrite-ssl.conf /etc/apache2/sites-available
a2enmod ssl
a2ensite default-ssl
a2enmod rewrite
a2ensite rewrite-ssl
a2dissite 000-default

installDate=$(date)
cp /tmp/roadapplepi.sql /tmp/roadapplepi-configd.sql
echo "INSERT INTO env (name, value) VALUES (\"rapVersion\", \"$softwareVersion\"), (\"installDate\", \"$installDate\");" >> /tmp/roadapplepi-configd.sql
mysql < /tmp/roadapplepi-configd.sql

echo "Done! Please reboot your Raspberry Pi now"
